# docker/dask/calcite_orchestrator_pyblazing/supervisord/dask-scheduler.conf
[supervisord]
environment=NVIDIA_VISIBLE_DEVICES=0

[program:scheduler]
command=/usr/bin/run_scheduler.sh
autostart=true
