# docker/dask/calcite_orchestrator_pyblazing/supervisord/dask-scheduler.conf
[program:scheduler]
command=/usr/bin/run_scheduler.sh
autostart=true
