# docker/dask/calcite_orchestrator/supervisord/dask-scheduler.conf
[program:worker]
command=dask-scheduler --show
autostart=true
