# docker/conda/blazingsql-template/blazingsql/runtime/supervisord/conf.d/calcite.conf
[program:calcite]
command=java -jar ../bin/BlazingCalcite.jar --data_directory=../catalog
autostart=false
