# docker/conda/blazingsql-template/blazingsql/runtime/supervisord/conf.d/orchestrator.conf
[program:orchestrator]
command=../bin/blazingdb_orchestator_service 8890 localhost 8891 localhost 8892
autostart=false
