# docker/conda/blazingsql-template/blazingsql/runtime/supervisord/conf.d/ral.conf
[program:ral]
environment=LD_LIBRARY_PATH="../cudf/cpp/build/"
command=../bin/testing-libgdf
autostart=true