# docker/conda/blazingsql-template/blazingsql/runtime/supervisord/conf.d/ral.conf
[program:ral]
environment=LD_LIBRARY_PATH=../lib/
command=../bin/testing-libgdf localhost 8892
autostart=false
