# docker/conda/blazingsql-template/supervisord/conf.d/calcite.conf
[program:calcite]
command=java -jar /repo/supervisor-demo/BlazingCalcite.jar
;directory=/blazingsql
;exitcodes=143
autostart=true
;user=jupyter
