# docker/conda/blazingsql-template/supervisord/conf.d/calcite.conf
[program:calcite]
command=java -jar ../bin/BlazingCalcite.jar
autostart=true