# docker/conda/blazingsql-template/supervisord/conf.d/orchestrator.conf
[program:orchestrator]
command=/repo/supervisor-demo/blazingdb_orchestator_service
autostart=true
;user=jupyter
