# docker/conda/blazingsql-template/supervisord/conf.d/orchestrator.conf
[program:orchestrator]
command=../bin/blazingdb_orchestator_service
autostart=true
