# docker/conda/blazingsql-template/supervisord/conf.d/ral.conf
[program:ral]
environment=A="1",B="2"
command=../bin/testing-libgdf
autostart=true