# docker/conda/blazingsql-template/supervisord/conf.d/ral.conf
[program:ral]
command=/repo/supervisor-demo/testing-libgdf
autostart=true
;user=jupyter
