# docker/blazingsql/supervisord/blazing-orchestrator.conf
[program:blazing-orchestrator]
command=/home/jupyter/blazingdb_orchestator_service 8890 localhost 8891 localhost 8892
autostart=true
user=jupyter
