# docker/blazingsql/supervisord/blazing-calcite.conf
[program:blazing-calcite]
command=java -Djava.library.path=/home/jupyter/lib -jar /home/jupyter/BlazingCalcite.jar
directory=/blazingsql
exitcodes=143
autostart=true
user=jupyter
