# docker/blazingsql/supervisord/blazing-orchestrator.conf
[program:blazing-orchestrator]
command=/home/jupyter/blazingdb_orchestator_service
autostart=true
user=jupyter
