# docker/blazingsql/supervisord/blazing-orchestrator.conf
[program:blazing-orchestrator]
command=/usr/bin/blazingdb_orchestator_service 8890 127.0.0.1 8891 127.0.0.1 8892
autostart=true
user=jupyter
