# docker/blazingsql/supervisord/blazing-calcite-tcp.conf
[program:blazing-calcite]
command=java -jar /home/jupyter/BlazingCalcite.jar -p 8890
directory=/blazingsql
exitcodes=143
autostart=true
user=jupyter
