# docker/blazingsql/supervisord/blazing-calcite-tcp.conf
[program:blazing-calcite]
command=java -jar /home/jupyter/BlazingCalcite.jar
directory=/blazingsql
exitcodes=143
autostart=false
user=jupyter
