# docker/blazingsql/supervisord/blazing-orchestrator.conf
[program:blazing-orchestrator]
command=blazingdb_orchestator_service
autostart=true
user=jupyter
