# docker/blazingsql/supervisord/blazing-orchestrator-tcp.conf
[program:blazing-orchestrator]
command=/home/jupyter/blazingdb_orchestator_service 8889 127.0.0.1 8890
autostart=true
user=jupyter
