# docker/blazingsql/supervisord/jupyterlab.conf
[program:jupyterlab]
command=/usr/bin/run_jupyter.sh
autostart=true
