# docker/blazingsql/supervisord/jupyterlab.conf
[program:jupyterlab]
command=run_jupyter.sh
autostart=true
