# docker/blazingsql/supervisord/jupyterlab.conf
[program:blazing-ral]
command=source activate gdf && bash /rapids/bash utils/start_jupyter.sh
autostart=true
