# docker/colab/supervisor/blazing-calcite.conf
[program:blazing-calcite]
command=java -jar /usr/bin/BlazingCalcite.jar -p 8890 --data_directory=/blazingsql
directory=/blazingsql
exitcodes=143
autostart=true
