# docker/colab/supervisor/blazing-orchestrator.conf
[program:blazing-orchestrator]
command=/usr/bin/blazingdb_orchestator_service 9000 8889 127.0.0.1 8890
autostart=true
