# docker/colab/supervisor/blazing-calcite.conf
[program:blazing-calcite]
command=java -jar /usr/bin/BlazingCalcite.jar
directory=/blazingsql
exitcodes=143
autostart=true
